# docker/crontab
* * * * * /usr/local/bin/php /var/www/html/artisan schedule:run
* * * * * /usr/local/bin/php /var/www/html/artisan queue:work --tries=3 --once --quiet